# zigzag-tools/zigzag-autoupdate.service
[Unit]
Description=Autoupdate the system

[Service]
Type=oneshot
ExecStart=/usr/bin/zypper dup -y
